# server-config/kamailio.service.d/override.conf
# Run in foreground (-D) so systemd keeps the service active (no PID file / forking)
[Service]
Type=simple
ExecStart=
ExecStart=/usr/sbin/kamailio -D -f /etc/kamailio/kamailio.cfg -m 64 -M 8
